# hxp20/pfoten/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::once:-setuidgid 1 sh
